# Honeybox/thug.conf
description     "Thug"
start on startup
task
exec /opt/thug/feedthug.sh